# obs-service-wsl-handler/systemd/obsidian-git-backup.service
[Unit]
Description=Backup Obsidian vault to GitHub
Wants=network-online.target
After=network-online.target obsidian-sync.service

[Service]
Type=oneshot
WorkingDirectory=%h
Environment=BOOTSTRAP_FORCE_PUSH=yes
ExecStart=%h/.local/bin/obsidian_git_backup.sh
